# tools/procServMgr/procServMgr.cron
Add the following lines to /etc/crontab.  This will ensure that softIOCs are up
and running.

# Added <date> by <author>
5,20,35,50 * * * * epics /usr/local/clas12/hps_epics/tools/bin/procServMgr check >> /usr/local/clas12/hps_epics/apps/iocBoot/softioc_`hostname`.log 2>&1
